# kodi-backup.service
[Unit]
Description=kodi-backup daemon

[Service]
ExecStart=/usr/local/sbin/kodi-backup
User=root
